# bin/99-lid.rules
# Trigger display blanking updates when the lid state or AC state changes
ACTION=="change", SUBSYSTEM=="button", KERNEL=="lid*", RUN+="/usr/bin/systemd-run --no-block /opt/pve-sleep/bin/lid-handler.sh"
ACTION=="change", SUBSYSTEM=="power_supply", ATTR{type}=="Mains", RUN+="/usr/bin/systemd-run --no-block /opt/pve-sleep/bin/lid-handler.sh"
