# bin/pve-lid-handler.service
[Unit]
Description=PVE Sleep lid display handler
ConditionPathExists=/opt/pve-sleep/bin/lid-handler.sh

[Service]
Type=oneshot
ExecStart=/opt/pve-sleep/bin/lid-handler.sh
